# supervisor/conf.d/serf-agent.conf
# supervisor file for serf

[program:serf-agent]
command=/files/serf-start.sh
numprocs=1
autostart=true
autorestart=true
